# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=/usr/local/sbin/php-fpm

[program:laravel-migrate]
command=/bin/sh -c "php /var/www/artisan migrate --force"
autostart=true
autorestart=true
